# systemd/systemd_files/git-update.service
[Unit]
Description=Update webdev-experiments repo

[Service]
Type=oneshot
WorkingDirectory=/srv/webdev
ExecStart=/usr/bin/git pull --ff-only
ExecStartPost=/bin/systemctl daemon-reload
