# systemd/systemd_files/git-update.timer
[Unit]
Description=Run git-update.service every 1 minute drink water

[Timer]
OnBootSec=60
OnUnitActiveSec=60
Persistent=true

[Install]
WantedBy=timers.target
